# 99-udev2dbus.rules
ACTION=="add",    RUN+="/bin/bash /usr/local/bin/udev2dbus.sh"
ACTION=="remove", RUN+="/bin/bash /usr/local/bin/udev2dbus.sh"

